# debian/n2rxd.prerm
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

invoke-rc.d n2rxd stop || true
#DEBHELPER#

exit 0


